# jackery_api.service
[Unit]
Description=Jackery API Data Collector
After=network-online.target
Wants=network-online.target

[Service]
# ★ 自分のユーザー名に変更してください（例: pi, tetsu など）
User=pi
# ★ スクリプトの実際のパスに変更してください
WorkingDirectory=/home/pi/Jackery
ExecStart=/usr/bin/python3 /home/pi/Jackery/jackery_api.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
